# wg-routing.service
[Unit]
After=network.target

[Service]
ExecStart=/root/wg-up.sh

[Install]
WantedBy=default.target